# hello.service
[Unit]
Description=My Hello world service

[Service]
Type=oneshot
ExecStart=/bin/echo "Hello, world!"

[Install]
WantedBy=multi-user.target
